# packages/lib/systemd/system/synse-server.service
[Unit]
Description=Synse Server
After=network-online.target
Wants=network.target
Requires=docker.service

[Service]
ExecStart=/usr/bin/docker run -v /var/log/synse:/logs -v /etc/synse/config.json:/synse/override/config.json -p 5000:5000 --name synse-server vaporio/synse-server
ExecStop=/usr/bin/docker rm -f synse-server
Restart=on-abnormal
RestartSec=20

[Install]
WantedBy=multi-user.target
