# packages/synse-server.systemd.postinst
#!/bin/sh
set -e
systemctl enable synse-server
